# config/smb-webadmin.service
[Unit]
Description=SMB WebAdmin Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/smb-webadmin/backend
ExecStart=/opt/smb-webadmin/backend/venv/bin/python /opt/smb-webadmin/backend/app.py
User=www-data
Group=www-data
Restart=on-failure

[Install]
WantedBy=multi-user.target
